# deploy/systemd/caller-worker.service
[Unit]
Description=TypeScript caller worker
After=network-online.target
Wants=network-online.target

[Service]
User=iii
WorkingDirectory=/opt/devops-assignment/quickstart/workers/caller-worker
Environment=III_URL=ws://10.10.0.10:49134
ExecStart=/usr/bin/node dist/worker.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
