# deploy/systemd/iii-engine.service
[Unit]
Description=iii engine and HTTP trigger worker
After=network-online.target
Wants=network-online.target

[Service]
User=iii
WorkingDirectory=/opt/iii
ExecStart=/usr/local/bin/iii --config /opt/iii/iii-config.yaml
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
