# deploy/systemd/inference-worker.service
[Unit]
Description=Python inference worker
After=network-online.target
Wants=network-online.target

[Service]
User=iii
WorkingDirectory=/opt/devops-assignment/quickstart/workers/inference-worker
Environment=III_URL=ws://10.10.0.10:49134
Environment=HF_HOME=/opt/iii/huggingface
ExecStart=/opt/devops-assignment/quickstart/workers/inference-worker/.venv/bin/python inference_worker.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
